# provisioning/roles/backup/templates/backup.service.j2
[Unit]
Description=run {{ backup_name }} backup

[Service]
Type=simple
ExecStart=/usr/local/sbin/{{ backup_name }}_backup.sh

[Install]
WantedBy=multi-user.target
