# provisioning/roles/backup/templates/backup.timer.j2
[Unit]
Description=run {{ backup_name }} backup on boot

[Timer]
OnBootSec={{ backup_time_delay }}
Persistent=true

[Install]
WantedBy=timers.target
