# provisioning/roles/backup/templates/test.timer.j2
[Unit]
Description=run {{ backup_name }} backup test every month

[Timer]
OnCalendar=Mon *-*-1..7 21:00:00
Persistent=true

[Install]
WantedBy=timers.target
