# roles/docker/templates/dockerhubmirror.service.j2
[Unit]
Description=dockerhubmirror
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/bin/bash /usr/share/systemd/dockerhubmirror.sh

[Install]
WantedBy=multi-user.target
