# provisioning/roles/mysql/templates/mysql.service.j2
[Unit]
Description=mysql
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/bin/bash /usr/share/systemd/mysqld.sh

[Install]
WantedBy=multi-user.target

