# provisioning/roles/nexus/templates/nexus.service.j2
[Unit]
Description=nexus
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/bin/bash /usr/share/systemd/nexus.sh

[Install]
WantedBy=multi-user.target

