# provisioning/roles/watchtower/templates/watchtower.service.j2
[Unit]
Description=watchtower
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/bin/bash /usr/share/systemd/watchtower.sh

[Install]
WantedBy=multi-user.target

