# roles/hints/templates/hintsd.service.j2
[Unit]
Description=Hints daemon
[Service]
Type=simple
ExecStart=/home/{{ ansible_user }}/.local/bin/hintsd
Restart=always
[Install]
WantedBy=default.target
