# roles/opendata/templates/opendata-download.timer.j2
[Unit]
Description=run {{ opendata_download_name }} opendata download on boot

[Timer]
OnBootSec={{ opendata_download_time_delay }}
Persistent=true

[Install]
WantedBy=timers.target
