# provisioning/roles/pgadmin/templates/pgadmin.service.j2
[Unit]
Description=pgadmin
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/bin/bash /usr/share/systemd/pgadmin.sh

[Install]
WantedBy=multi-user.target




