# roles/chezmoi/templates/chezmoi-diff.service.j2
[Unit]
Description=run chezmoi-diff

[Service]
Type=simple
User={{ ansible_user }}
Group={{ ansible_user }}
UMask=0002
ExecStart=/usr/local/sbin/chezmoi-diff.sh

[Install]
WantedBy=multi-user.target
