# roles/chezmoi/templates/chezmoi-diff.timer.j2
[Unit]
Description=run chezmoi-diff on boot

[Timer]
OnBootSec=7200
Persistent=true

[Install]
WantedBy=timers.target
